# deploy/kasirtoko.service
[Unit]
Description=KasirToko — Flask POS App
After=network.target

[Service]
# Ganti 'ubuntu' dengan username Linux kamu
User=ubuntu
Group=ubuntu

WorkingDirectory=/home/ubuntu/kasirtoko
Environment="PATH=/home/ubuntu/kasirtoko/venv/bin"

# Gunicorn: 2 worker cukup untuk pemakaian rumahan
ExecStart=/home/ubuntu/kasirtoko/venv/bin/gunicorn \
    --workers 2 \
    --bind 127.0.0.1:5000 \
    --timeout 60 \
    --access-logfile /home/ubuntu/kasirtoko/logs/access.log \
    --error-logfile  /home/ubuntu/kasirtoko/logs/error.log \
    app:app

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
